# ansible/templates/systemd/user/redshift.service.d/override.conf
[Unit]
Wants=geoclue-agent.service
After=geoclue-agent.service
